# systemd/respondedorbot-podman-prune.service
[Unit]
Description=Prune unused Podman images for the respondedorbot user

[Service]
Type=oneshot
ExecStart=/usr/bin/podman image prune -f

[Install]
WantedBy=default.target
